# systemd/oxy.service
[Unit]
Description=Oxygen workspace server (SPA + agents on :3000)
# Critical: oxy start manages a postgres container via docker, so it must
# start AFTER docker is up. Without Requires/After=docker.service the unit
# races docker on reboot and crashes bringing up the container.
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory={{PROJECT_ROOT}}
EnvironmentFile=/etc/environment
# oxy lives in the install user's ~/.local/bin (script 03). Do NOT use %h here:
# systemd's %h resolves to the *service manager's* home (/root for system
# units), NOT the User= home, so %h/.local/bin/oxy → /root/.local/bin/oxy and
# fails 203/EXEC. Script 08 substitutes {{HOME_DIR}} with the install user's
# real home instead.
ExecStart={{HOME_DIR}}/.local/bin/oxy start --local
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
